# _ops/haminer-test/mkosi.extra/etc/systemd/system/haminer-dummy-backend.path
# SPDX-FileCopyrightText: 2024 M. Shulhan <ms@kilabit.info>
#
# SPDX-License-Identifier: GPL-3.0-or-later

[Unit]
Description="Watch haminer-dummy-backend"

[Path]
PathChanged=/data/haminer/bin/haminer-dummy-backend
Unit=systemctl-restart@%p.service

[Install]
WantedBy=default.target
