# _ops/haminer-test/mkosi.extra/etc/systemd/system/haminer-dummy-backend.path
[Unit]
Description="Watch haminer-dummy-backend"

[Path]
PathChanged=/data/haminer/bin/haminer-dummy-backend
Unit=systemctl-restart@%p.service

[Install]
WantedBy=default.target
