# _ops/haminer-test/mkosi.extra/etc/systemd/system/systemctl-restart@.service
[Unit]
Description=systemctl-restart@%i

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart %i
